# duplicacy-daily-job.service
[Unit]
Description=Daily duplicacy backup script

[Service]
ExecStart=/home/t/duplicacy-script/duplicacy-backup
Type=oneshot
